# vm/host-services/open-claw/systemd/openclaw-gateway.service.example
# This file is auto-generated by the openclaw CLI (openclaw gateway install).
# It is stored here as a reference template only — do NOT symlink it.
# After running `openclaw gateway install`, the actual service file will be at:
#   ~/.config/systemd/user/openclaw-gateway.service

[Unit]
Description=OpenClaw Gateway (v2026.2.6-3)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart="/usr/bin/node" "%h/.npm-global/lib/node_modules/openclaw/dist/index.js" gateway --port 18789
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=%h
Environment="PATH=%h/.local/bin:%h/.npm-global/bin:%h/bin:%h/.nvm/current/bin:%h/.fnm/current/bin:%h/.volta/bin:%h/.asdf/shims:%h/.local/share/pnpm:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_GATEWAY_TOKEN=<SET_IN_SECRETS_ENV>
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.2.6-3

[Install]
WantedBy=default.target
