# vm/host-services/open-claw/systemd/openclaw-proxy.service.example
# This is a reference template — copy and adjust before deploying.
# To install:
#   cp openclaw-proxy.service.example ~/.config/systemd/user/openclaw-proxy.service
#   # Edit S_LLM_API_URL to match your environment
#   systemctl --user daemon-reload
#   systemctl --user enable --now openclaw-proxy.service

[Unit]
Description=OpenClaw OpenAI-compat Proxy
Before=openclaw-gateway.service

[Service]
Type=simple
ExecStart=/usr/bin/node %h/.openclaw/openai-compat-proxy.js
Environment=S_LLM_API_URL=<SET_YOUR_LLM_API_URL>
Environment=PROXY_PORT=19999
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
